# package/base-files/files/etc/init.d/done
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=95
boot() {
	mount_root done
	rm -f /sysupgrade.tgz && sync
<<<<<<< HEAD
	rm -f /tmp/sysupgrade.tar && sync
=======
>>>>>>> 712839d4c6 (Removed unwanted submodules from index)

	# process user commands
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}

	# set leds to normal state
	. /etc/diag.sh
	set_state done
}
